# docker/servlet-engine/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:appserver]
directory = /opt/appserver
command = bin/php -dappserver.php_sapi=appserver -f /opt/appserver/server.php
user = root
autostart = true